# server/backend.ini
[program:backend]
directory=/home/USERNAME/webserver
command=npm run start
autostart=true
autorestart=true
startsecs=20                    
startretries=5                     
stopwaitsecs=10
environment=NODE_ENV=production